# bastion-firewall.service
[Unit]
Description=Bastion Firewall Daemon
After=network.target

[Service]
Type=notify
ExecStart=/usr/bin/python3 /usr/bin/bastion-daemon
Restart=on-failure
RestartSec=5s
WatchdogSec=30s
User=root
Group=root

[Install]
WantedBy=multi-user.target
