# bastion-rs/bastion-daemon.service
[Unit]
Description=Bastion Application Firewall Daemon
After=network-online.target
Wants=network-online.target
Documentation=https://github.com/shipdocs/bastion-firewall

[Service]
Type=simple
ExecStartPre=/usr/sbin/iptables -I OUTPUT 1 -m state --state NEW -j NFQUEUE --queue-num 1 --queue-bypass
ExecStart=/usr/bin/bastion-daemon
ExecStopPost=/usr/sbin/iptables -D OUTPUT -m state --state NEW -j NFQUEUE --queue-num 1 --queue-bypass
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bastion-daemon

# Directories
RuntimeDirectory=bastion
ConfigurationDirectory=bastion

# Security hardening
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_PTRACE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_PTRACE

[Install]
WantedBy=multi-user.target
