# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Create bastion user/group if they don't exist
if ! getent group bastion >/dev/null; then
    groupadd -r bastion
fi
if ! getent passwd bastion >/dev/null; then
    useradd -r -g bastion -d /var/lib/bastion -s /usr/sbin/nologin -c "Bastion Firewall User" bastion
fi

# Set up configuration directory
mkdir -p /etc/bastion
if [ ! -f /etc/bastion/config.json ]; then
    if [ -f /usr/share/doc/bastion-firewall/config.json.example ]; then
        cp /usr/share/doc/bastion-firewall/config.json.example /etc/bastion/config.json
    else
        echo '{"log_level": "INFO", "rules_file": "/etc/bastion/rules.json"}' > /etc/bastion/config.json
    fi
    chmod 644 /etc/bastion/config.json
fi

# Set up rules file
if [ ! -f /etc/bastion/rules.json ]; then
    echo '{"applications": {}, "services": {}}' > /etc/bastion/rules.json
    chmod 644 /etc/bastion/rules.json
fi

# Set permissions for logs
mkdir -p /var/log/bastion
chown -R bastion:bastion /var/log/bastion
chmod 750 /var/log/bastion

# Install logrotate configuration
if [ -f /usr/share/doc/bastion-firewall/bastion-firewall.logrotate ]; then
    cp /usr/share/doc/bastion-firewall/bastion-firewall.logrotate /etc/logrotate.d/bastion-firewall
fi

# Ensure log file has correct permissions if it exists
touch /var/log/bastion-daemon.log
chgrp bastion /var/log/bastion-daemon.log
chmod 644 /var/log/bastion-daemon.log

# Verify daemon binary
if [ ! -x /usr/bin/bastion-daemon ]; then
    echo "ERROR: Daemon binary not found"
    exit 1
fi

# Warn if eBPF unavailable
if [ ! -f /usr/share/bastion-firewall/bastion-ebpf.o ] || [ ! -f /sys/kernel/btf/vmlinux ]; then
    echo "WARNING: eBPF not available. Daemon will use /proc fallback."
fi

mkdir -p /var/run/bastion

# Reload systemd
systemctl daemon-reload

# Enable service
systemctl enable bastion-firewall

# Try starting service and verify it succeeds
echo "Starting bastion-firewall service..."
if systemctl start bastion-firewall; then
    # Wait a moment and verify it's still running
    sleep 2
    if ! systemctl is-active --quiet bastion-firewall; then
        echo "✗ ERROR: Service failed to start after 2 seconds. Installation failed."
        echo "Check logs with: journalctl -u bastion-firewall -n 50"
        # Stop service to ensure cleanup
        systemctl stop bastion-firewall || true
        exit 1
    fi
    echo "✓ bastion-firewall service started successfully"
else
    echo "✗ ERROR: Failed to start bastion-firewall service. Installation failed."
    exit 1
fi

# Update desktop database
update-desktop-database || true

# Autostart helpers
if [ -d /home ]; then
    for user_dir in /home/*; do
        if [ -d "$user_dir/.config/autostart" ]; then
            chown -R $(basename "$user_dir"):$(basename "$user_dir") "$user_dir/.config/autostart" 2>/dev/null || true
        fi
    done
fi
