# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Colors for output
RED='\033[0;31m'
GREEN='\033[0;32m'
NC='\033[0m'

print_success() { echo -e "${GREEN}✓ $1${NC}"; }
print_error() { echo -e "${RED}✗ $1${NC}"; }

# Create bastion user/group if they don't exist
if ! getent group bastion >/dev/null; then
    groupadd -r bastion
fi
if ! getent passwd bastion >/dev/null; then
    useradd -r -g bastion -d /var/lib/bastion -s /usr/sbin/nologin -c "Bastion Firewall User" bastion
fi

# Set up configuration directory
mkdir -p /etc/bastion
if [ ! -f /etc/bastion/config.json ]; then
    if [ -f /usr/share/doc/bastion-firewall/config.json.example ]; then
        cp /usr/share/doc/bastion-firewall/config.json.example /etc/bastion/config.json
    else
        echo '{"log_level": "INFO", "rules_file": "/etc/bastion/rules.json"}' > /etc/bastion/config.json
    fi
    chmod 644 /etc/bastion/config.json
fi

# Set up rules file
if [ ! -f /etc/bastion/rules.json ]; then
    echo '{"applications": {}, "services": {}}' > /etc/bastion/rules.json
    chmod 644 /etc/bastion/rules.json
fi

# Set permissions for logs
mkdir -p /var/log/bastion
chown -R bastion:bastion /var/log/bastion
chmod 750 /var/log/bastion

# Install logrotate configuration
if [ -f /usr/share/doc/bastion-firewall/bastion-firewall.logrotate ]; then
    cp /usr/share/doc/bastion-firewall/bastion-firewall.logrotate /etc/logrotate.d/bastion-firewall
fi

# Ensure log file has correct permissions if it exists
if [ -f /var/log/bastion-daemon.log ]; then
    chgrp bastion /var/log/bastion-daemon.log
    chmod 640 /var/log/bastion-daemon.log
fi

# Reload systemd
# Reload systemd
systemctl daemon-reload
# Always enable and start/restart firewall
systemctl enable bastion-firewall
systemctl restart bastion-firewall

# Update desktop database
update-desktop-database || true

# Autostart helpers
if [ -d /home ]; then
    for user_dir in /home/*; do
        if [ -d "$user_dir/.config/autostart" ]; then
            chown -R $(basename "$user_dir"):$(basename "$user_dir") "$user_dir/.config/autostart" 2>/dev/null || true
        fi
    done
fi
