# debian/DEBIAN/postinst
#!/bin/bash
set -e

# Colors for output
RED='\033[0;31m'
GREEN='\033[0;32m'
NC='\033[0m'

print_success() { echo -e "${GREEN}✓ $1${NC}"; }
print_error() { echo -e "${RED}✗ $1${NC}"; }

# SECURITY: Create bastion user/group with proper permissions
if ! getent group bastion >/dev/null; then
    groupadd -r bastion
    print_success "Created bastion group"
fi
if ! getent passwd bastion >/dev/null; then
    useradd -r -g bastion -d /var/lib/bastion -s /usr/sbin/nologin -c "Bastion Firewall User" bastion
    print_success "Created bastion user"
fi

# Add first regular user to bastion group for convenience
FIRST_USER=$(getent passwd 1000 | cut -d: -f1)
if [ -n "$FIRST_USER" ]; then
    if ! groups "$FIRST_USER" | grep -q bastion; then
        usermod -a -G bastion "$FIRST_USER"
        print_success "Added user $FIRST_USER to bastion group (requires logout)"
    fi
fi

# Set up configuration directory
mkdir -p /etc/bastion
if [ ! -f /etc/bastion/config.json ]; then
    if [ -f /usr/share/doc/bastion-firewall/config.json.example ]; then
        cp /usr/share/doc/bastion-firewall/config.json.example /etc/bastion/config.json
    else
        echo '{"log_level": "INFO", "rules_file": "/etc/bastion/rules.json"}' > /etc/bastion/config.json
    fi
    # SECURITY FIX: Restrict config permissions to root only
    chmod 600 /etc/bastion/config.json
    print_success "Created config with secure permissions (600)"
fi

# Set up rules file
if [ ! -f /etc/bastion/rules.json ]; then
    echo '{"applications": {}, "services": {}}' > /etc/bastion/rules.json
    # SECURITY FIX: Restrict rules permissions to root only
    chmod 600 /etc/bastion/rules.json
    print_success "Created rules with secure permissions (600)"
fi

# Set permissions for logs
mkdir -p /var/log/bastion
chown -R root:bastion /var/log/bastion
chmod 750 /var/log/bastion
print_success "Set up log directory with restricted permissions"

# Reload systemd
systemctl daemon-reload
if systemctl is-enabled bastion-firewall >/dev/null 2>&1; then
    systemctl enable bastion-firewall
fi

# Update desktop database
update-desktop-database || true

# Autostart helpers
if [ -d /home ]; then
    for user_dir in /home/*; do
        if [ -d "$user_dir/.config/autostart" ]; then
            chown -R $(basename "$user_dir"):$(basename "$user_dir") "$user_dir/.config/autostart" 2>/dev/null || true
        fi
    done
fi

echo ""
echo "======================================================================"
echo "  Bastion Firewall Installation Complete"
echo "======================================================================"
echo ""
echo "SECURITY NOTES:"
echo "  - Config/rules files are now root-only (600 permissions)"
echo "  - Socket will be restricted to bastion group (660 permissions)"
echo "  - First user added to bastion group automatically"
echo ""
echo "To add more users: sudo usermod -a -G bastion USERNAME"
echo "To start: bastion-gui or from application menu"
echo ""
echo "Documentation: /usr/share/doc/bastion-firewall/"
echo "======================================================================"
echo ""
