# debian/DEBIAN/prerm
#!/bin/bash
set -e
systemctl stop bastion-firewall || true
systemctl disable bastion-firewall || true
pkill -f bastion-gui || true
pkill -f bastion-control-panel || true
