# debian/DEBIAN/prerm
#!/bin/bash
#
# Pre-removal script for Bastion Firewall
# This runs BEFORE package files are removed
#

set -e

echo ""
echo "============================================================"
echo "ðŸ° Bastion Firewall - Pre-Removal"
echo "============================================================"
echo ""

# Stop and disable the service
if systemctl is-active --quiet bastion-firewall 2>/dev/null; then
    echo "Stopping bastion-firewall service..."
    systemctl stop bastion-firewall 2>/dev/null || true
    echo "âœ“ Service stopped"
fi

if systemctl is-enabled --quiet bastion-firewall 2>/dev/null; then
    echo "Disabling bastion-firewall service..."
    systemctl disable bastion-firewall 2>/dev/null || true
    echo "âœ“ Service disabled"
fi

# Kill all Bastion processes (daemon and GUIs)
echo "Terminating all Bastion processes..."
pkill -f bastion-daemon 2>/dev/null || true
pkill -f bastion-gui 2>/dev/null || true
pkill -f bastion-control-panel 2>/dev/null || true
sleep 1
echo "âœ“ Processes terminated"

# Remove iptables NFQUEUE rules
echo "Removing iptables NFQUEUE rules..."
iptables -D OUTPUT -m state --state NEW -j NFQUEUE --queue-num 1 2>/dev/null || true
echo "âœ“ Firewall rules removed"

# Remove socket file
if [ -S "/tmp/bastion-daemon.sock" ]; then
    echo "Removing socket file..."
    rm -f /tmp/bastion-daemon.sock
    echo "âœ“ Socket removed"
fi

echo ""
echo "============================================================"
echo "âœ“ Bastion Firewall stopped and cleaned up"
echo "============================================================"
echo ""
echo "âš  IMPORTANT: Your UFW outbound policy may still be set to DENY."
echo ""
echo "To restore normal outbound access, run:"
echo "  sudo ufw default allow outgoing"
echo "  sudo ufw reload"
echo ""

exit 0

