# debian/DEBIAN/prerm
#!/bin/bash
#
# Pre-removal script for Douane Firewall
# This runs BEFORE package files are removed
#

set -e

echo ""
echo "============================================================"
echo "Douane Firewall - Pre-Removal"
echo "============================================================"
echo ""

# Stop and disable the service
if systemctl is-active --quiet douane-firewall 2>/dev/null; then
    echo "Stopping douane-firewall service..."
    systemctl stop douane-firewall 2>/dev/null || true
    echo "✓ Service stopped"
fi

if systemctl is-enabled --quiet douane-firewall 2>/dev/null; then
    echo "Disabling douane-firewall service..."
    systemctl disable douane-firewall 2>/dev/null || true
    echo "✓ Service disabled"
fi

# Kill all Douane processes (daemon and GUIs)
echo "Terminating all Douane processes..."
pkill -f douane-daemon 2>/dev/null || true
pkill -f douane-gui-client 2>/dev/null || true
pkill -f douane-control-panel 2>/dev/null || true
pkill -f douane_gui_improved 2>/dev/null || true
sleep 1
echo "✓ Processes terminated"

# Remove iptables NFQUEUE rules
echo "Removing iptables NFQUEUE rules..."
iptables -D OUTPUT -m state --state NEW -j NFQUEUE --queue-num 1 2>/dev/null || true
echo "✓ Firewall rules removed"

# Remove socket file
if [ -S "/var/run/douane.sock" ]; then
    echo "Removing socket file..."
    rm -f /var/run/douane.sock
    echo "✓ Socket removed"
fi

echo ""
echo "============================================================"
echo "✓ Douane Firewall stopped and cleaned up"
echo "============================================================"
echo ""
echo "⚠ IMPORTANT: Your UFW outbound policy may still be set to DENY."
echo ""
echo "To restore normal outbound access, run:"
echo "  sudo ufw default allow outgoing"
echo "  sudo ufw reload"
echo ""

exit 0

