# debian/DEBIAN/prerm
#!/bin/bash
#
# Pre-removal script for Bastion Firewall
# This runs BEFORE package files are removed
#

set -e

echo ""
echo "============================================================"
echo "ðŸ° Bastion Firewall - Pre-Removal"
echo "============================================================"
echo ""

# Stop and disable the service (both Bastion and old Douane)
if systemctl is-active --quiet bastion-firewall 2>/dev/null; then
    echo "Stopping bastion-firewall service..."
    systemctl stop bastion-firewall 2>/dev/null || true
    echo "âœ“ Bastion service stopped"
fi

if systemctl is-active --quiet douane 2>/dev/null; then
    echo "Stopping old douane service..."
    systemctl stop douane 2>/dev/null || true
    echo "âœ“ Old Douane service stopped"
fi

if systemctl is-enabled --quiet bastion-firewall 2>/dev/null; then
    echo "Disabling bastion-firewall service..."
    systemctl disable bastion-firewall 2>/dev/null || true
    echo "âœ“ Service disabled"
fi

# Kill ALL Bastion processes (daemon and GUIs)
echo "Terminating all Bastion processes..."
pkill -f bastion-daemon 2>/dev/null || true
pkill -f bastion-gui 2>/dev/null || true
pkill -f bastion-control-panel 2>/dev/null || true
pkill -f bastion-tray 2>/dev/null || true
sleep 0.5

# Force kill if still running
pkill -9 -f bastion-daemon 2>/dev/null || true
pkill -9 -f bastion-gui 2>/dev/null || true
echo "âœ“ Bastion processes terminated"

# Kill ALL old Douane processes
echo "Terminating all old Douane processes..."
pkill -f douane-daemon 2>/dev/null || true
pkill -f douane-gui 2>/dev/null || true
pkill -f douane-control-panel 2>/dev/null || true
pkill -f douane-tray 2>/dev/null || true
sleep 0.5

# Force kill if still running
pkill -9 -f douane-daemon 2>/dev/null || true
pkill -9 -f douane-gui 2>/dev/null || true
echo "âœ“ Old Douane processes terminated"

# Remove iptables NFQUEUE rules (all possible queue numbers)
echo "Removing iptables NFQUEUE rules..."
iptables -D OUTPUT -j NFQUEUE --queue-num 0 2>/dev/null || true
iptables -D OUTPUT -j NFQUEUE --queue-num 1 2>/dev/null || true
iptables -D OUTPUT -m state --state NEW -j NFQUEUE --queue-num 0 2>/dev/null || true
iptables -D OUTPUT -m state --state NEW -j NFQUEUE --queue-num 1 2>/dev/null || true
echo "âœ“ Firewall rules removed"

# Remove socket files (both Bastion and old Douane)
echo "Removing socket files..."
rm -f /tmp/bastion-daemon.sock 2>/dev/null || true
rm -f /var/run/douane.sock 2>/dev/null || true
rm -f /tmp/douane-daemon.sock 2>/dev/null || true
echo "âœ“ Socket files removed"

echo ""
echo "============================================================"
echo "âœ“ Bastion Firewall stopped and cleaned up"
echo "============================================================"
echo ""
echo "âš  IMPORTANT: Your UFW outbound policy may still be set to DENY."
echo ""
echo "To restore normal outbound access, run:"
echo "  sudo ufw default allow outgoing"
echo "  sudo ufw reload"
echo ""

exit 0

